# deploy/gcp-lab/needletail-mesh.service
[Unit]
Description=Needletail AV Mesh qualification node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/needletail/node.env
Environment=RUST_LOG=av_mesh=info,av_web_service=warn
Environment=NEEDLETAIL_MISSION_CONTROL_DIST=/opt/needletail/mission-control
ExecStart=/usr/local/bin/needletail-av-mesh-run
Restart=always
RestartSec=2
NoNewPrivileges=yes
PrivateTmp=yes
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
