# deploy/gcp-lab/needletail-media.service
[Unit]
Description=Needletail synthetic realtime qualification source
After=needletail-contrib.service
Requires=needletail-contrib.service

[Service]
Type=simple
ExecStart=/usr/bin/ffmpeg -hide_banner -loglevel warning -re -f lavfi -i testsrc2=size=960x540:rate=30 -f lavfi -i sine=frequency=997:sample_rate=48000 -c:v libx264 -preset ultrafast -tune zerolatency -g 30 -keyint_min 30 -sc_threshold 0 -b:v 1200k -c:a aac -b:a 96k -ar 48000 -f flv rtmp://127.0.0.1:19350/live/gcp-qualification
Restart=always
RestartSec=2
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
