# deploy/gcp-lab/needletail-rist-source-lori-4k.service
[Unit]
Description=Needletail LORI 4K RIST source
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/mnt/needletail-media

[Service]
Type=simple
User=jamie
Group=jamie
KillMode=control-group
ExecStart=/bin/bash -lc 'set -o pipefail; exec /usr/bin/ffmpeg -hide_banner -loglevel warning -re -stream_loop -1 -i /mnt/needletail-media/lori_4k_no_grain_4k25_ll_capped10.mp4 -map 0:v:0 -map 0:a:0 -c copy -muxdelay 0 -muxpreload 0 -f mpegts - | /usr/local/bin/rist-send --profile main --flow-id 0x11223344 --chunk-bytes 1316 --history-packets 65536 --final-repair-ms 1000 10.84.10.5:27000'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
